# debian/lcd4linux-ax206.service
[Unit]
Description=LCD4Linux display daemon for AX206 USB displays
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/share/lcd4linux
ExecStartPre=/bin/sh -c 'test -f /etc/lcd4linux/lcd4linux.conf || { echo "No config found. Run: sudo lcd4linux-theme apply SimpleBlue"; exit 1; }'
ExecStart=/usr/bin/lcd4linux -F -f /etc/lcd4linux/lcd4linux.conf
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
